# systemd/uki-btrfs-snapshots.service
[Unit]
Description=Clone Unified Kernel Images per btrfs snapshot
Documentation=https://github.com/jmylchreest/refind-btrfs-snapshots
After=local-fs.target
Before=snapper-boot.service
ConditionPathExists=/usr/bin/uki-btrfs-snapshots
ConditionPathExists=/etc/uki-btrfs-snapshots.yaml

[Service]
Type=oneshot
# Settle window: a kernel rebuild or snapshot rotation produces a burst of
# writes. Sleep briefly so they coalesce; because this is Type=oneshot, triggers
# arriving during the sleep/run collapse into at most one re-run.
ExecStartPre=/usr/bin/sleep 5
ExecStart=/usr/bin/uki-btrfs-snapshots --config /etc/uki-btrfs-snapshots.yaml generate -y
User=root
StandardOutput=journal
StandardError=journal
TimeoutStartSec=300
SyslogIdentifier=uki-btrfs-snapshots

[Install]
WantedBy=multi-user.target
